# packaging/deb/postinst
#!/bin/sh
set -e
case "$1" in
    configure)
        if ! getent passwd merovingian > /dev/null 2>&1; then
            adduser --system --group --no-create-home \
                    --home /var/lib/merovingian \
                    --shell /usr/sbin/nologin \
                    --gecos "Merovingian homeserver" \
                    merovingian
        fi
        install -d -o merovingian -g merovingian -m 0750 /var/lib/merovingian
        install -d -o merovingian -g merovingian -m 0750 /var/log/merovingian

        # Generate a registration token on first install so the server can
        # start without manual setup. Never overwrite an existing token.
        TOKEN_FILE=/etc/merovingian/registration-token
        if [ ! -f "${TOKEN_FILE}" ]; then
            install -d -m 0755 /etc/merovingian
            openssl rand -base64 48 > "${TOKEN_FILE}"
            chmod 0640 "${TOKEN_FILE}"
            chown root:merovingian "${TOKEN_FILE}"
        fi
        ;;
esac
if [ -d /run/systemd/system ] && [ "$1" = "configure" ]; then
    systemctl --no-reload enable merovingian.service >/dev/null 2>&1 || true
fi
exit 0
